# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── Laravel Octane (FrankenPHP) — HTTP control plane ──────────
[program:octane]
command=php artisan octane:start --server=frankenphp --host=0.0.0.0 --port=8000 --admin-port=2019 --max-requests=500
directory=/app
autostart=true
autorestart=true
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ── Laravel queue worker — webhooks + jobs ────────────────────
[program:queue]
command=php artisan queue:work --sleep=1 --tries=3 --max-time=3600 --timeout=300
directory=/app
autostart=true
autorestart=true
stopwaitsecs=305
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

; ── Laravel scheduler — daily rollup prune, etc. ──────────────
[program:scheduler]
command=php artisan schedule:work
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30

; ── golang-queue — VOD transcode worker (ffmpeg → HLS → S3) ───
[program:go-queue]
command=/usr/local/bin/oxygen-queue
directory=/app
autostart=true
autorestart=true
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=40

; ── golang-live — RTMP ingest + live HLS ──────────────────────
[program:go-live]
command=/usr/local/bin/oxygen-live
directory=/app
autostart=true
autorestart=true
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=50
